# deploy/com.stbarnabas.social.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stbarnabas.social.watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/stbsocial/scripts/watch.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/stbsocial</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/stbsocial/cache/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/stbsocial/cache/launchd.err.log</string>
</dict>
</plist>
